# vuln_machine/crontab
*/5 * * * * root echo "root{$(openssl rand -hex 32)}" > /root/root.txt
*/5 * * * * ctf echo "user{$(openssl rand -hex 32)}" > /home/ctf/user.txt
